# supervisord.conf
[supervisord]
nodaemon=true
[program:flask]
directory = /opt/ssl-cert-exporter
command = /usr/local/bin/python ./app.py
autorestart = true

